# src/etc/udev/rules.d/98-wget-sync.rules
KERNEL=="eth*", RUN+="/usr/local/wget-sync/udev_program.sh" 
KERNEL=="wlan*", RUN+="/usr/local/wget-sync/udev_program.sh" 
KERNEL=="usb_host", RUN+="/usr/local/wget-sync/udev_usb_program.sh"
KERNEL=="usb_plug", RUN+="/usr/local/wget-sync/udev_usb_program.sh"

